# supervisord.conf
[supervisord]
nodaemon=true

[program:ssr]
stdout_logfile=/code/ssr.out
command=/code/bin/entrypoint-ssr

[program:web]
stdout_logfile=/code/web.out
command=/code/bin/entrypoint-web
